# puppet/modules/custom_nodejs/templates/nodejs_app.conf.erb
#nodejs_app
description "launch an app"
start on runlevel [2345]
stop on runlevel [!2345]
respawn
chdir /home/www-data/
exec su -s /bin/sh -c 'exec "$0" "$@"' www-data -- node ./nodejs-helloworld/app.js > nodejs-helloworld_log.log 2> nodejs-helloworld_err.log